# deploy/daphne.conf
[program:daphne]
command=/root/web3-wallet/env/bin/daphne -p 8000 project.asgi:application
directory=/root/web3-wallet/project
autostart=true
autorestart=true
stderr_logfile=/var/log/daphne.err.log
stdout_logfile=/var/log/daphne.out.log